# tld-server/etc-supervisord.d-supervisord.ini
[supervisord]
nodaemon=true

[program:mariadb]
command=mysqld_safe
process_name=/usr/libexec/mysqld
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
priority=1
startretries=15
autorestart=true

[program:pdns]
command=pdns_server --daemon=no
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
priority=999
startretries=15
autorestart=true
