# tld-server/etc-supervisord.d-supervisord.ini
[supervisord]
nodaemon=true

[program:mariadb]
command=mysqld_safe
process_name=/usr/libexec/mysqld
redirect_stderr=true

[program:pdns]
command=pdns_server --daemon=no
